# scripts/print-server.service.example
# Printer API — copy to /etc/systemd/system/print-server.service, edit paths, then:
#   sudo systemctl daemon-reload && sudo systemctl enable --now print-server.service

[Unit]
Description=Thermal printer API (Express)
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/johnpatrickperkins/Desktop/printerGUI
User=johnpatrickperkins
ExecStart=/home/johnpatrickperkins/.nvm/versions/node/v24.14.1/bin/npm run server
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
