# scripts/printer-gui.service.example
# Copy to /etc/systemd/system/printer-gui.service, set WorkingDirectory and User, then:
#   sudo systemctl daemon-reload && sudo systemctl enable --now printer-gui.service

[Unit]
Description=Printer GUI (Vite dev + print server)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/pi/printerGUI
User=pi
ExecStart=/usr/bin/npm run start:pi
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
