# scripts/vite-ui.service.example
# Vite dev UI — separate unit from the printer server. Copy and enable after print-server:
#   sudo systemctl daemon-reload && sudo systemctl enable --now vite-ui.service

[Unit]
Description=Vite dev server (printer GUI UI)
After=network-online.target print-server.service

[Service]
Type=simple
WorkingDirectory=/home/pi/printerGUI
User=pi
ExecStart=/usr/bin/npm run dev -- --host
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
